# root/etc/services.d/codi-md/run
#!/usr/bin/with-contenv bash

cd /app/codi-md || exit

#exec \
#	s6-setuidgid abc
