# root/etc/services.d/codi-md/run
#!/usr/bin/with-contenv bash

cd /opt/codimd || exit

#exec \
#	s6-setuidgid abc
